# deploy/systemd/recta-micro.timer
[Unit]
Description=Micro-nouvelle — tous les 3-4 jours (Lun,Jeu 18:00), langue tournante

[Timer]
OnCalendar=Mon,Thu *-*-* 18:00:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
