# deploy/systemd/recta-micro.timer
[Unit]
Description=Micro-nouvelle — quotidienne (18:00), langue tournante

[Timer]
OnCalendar=*-*-* 18:00:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
