# deploy/systemd/recta-publish.timer
[Unit]
Description=Feuilleton ROBOTARIIS — beat du jour (1×/j, 09:00)

[Timer]
OnCalendar=*-*-* 09:00:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
