# deploy/systemd/recta-publish.service
[Unit]
Description=Recta — publication quotidienne de l'Oraculum sur Facebook
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/olivier/DEV/Recta/scripts/broadcast.sh communique
